# p2pool-observer-api.service
# /etc/systemd/system/p2pool-observer-api.service

[Unit]
Description=P2Pool Salvium Observer API
Documentation=https://github.com/salvium/p2pool-salvium-observer
After=network-online.target valkey.service p2pool-observer-daemon.service
Wants=network-online.target
Requires=valkey.service

[Service]
Type=simple
User=p2pool
Group=p2pool
WorkingDirectory=/srv/p2pool-observer

# Clean API - reads from Redis and serves JSON endpoints
ExecStart=/srv/p2pool-observer/api \
  --redis=localhost:6379 \
  --listen=:9090

# Graceful shutdown
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

# Restart policy
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/srv/p2pool-observer

# Resource limits (API is very lightweight - just reads from Redis)
LimitNOFILE=65535
MemoryMax=256M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=p2pool-observer-api

[Install]
WantedBy=multi-user.target
