# p2pool-listener.service
# /etc/systemd/system/p2pool-listener.service

[Unit]
Description=P2Pool Listener - Salvium P2Pool Block Capture
After=network.target

[Service]
Type=simple
User=p2pool
Group=p2pool
WorkingDirectory=/srv/p2pool-listener

# Path to the p2plistener binary
ExecStart=/srv/p2pool-listener/p2plistener \
  --sidechain-config=/srv/p2pool-listener/sidechain_config.json \
  --redis=localhost:6379 \
  --salvium-host=core2.whiskymine.io \
  --salvium-port=19081 \
  --listen=0.0.0.0:38889

# Restart on failure
Restart=on-failure
RestartSec=10s

# Resource limits
MemoryMax=512M
TasksMax=256

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=p2pool-listener

[Install]
WantedBy=multi-user.target
